# compose/django/s6/celery/run
#!/bin/sh

exec cd /app && celery -A velo.taskapp worker -l INFO --autoscale=5,2
