# compose/django/s6/celery/run
#!/usr/bin/with-contenv sh
cd /app
exec celery -A velo.taskapp worker -l INFO --autoscale=5,2
