# compose/elk/s6/logstash/run
#!/usr/bin/with-contenv sh

exec /opt/logstash/bin/logstash -f /etc/logstash/logstash.json
