# compose/elk/s6/kibana/run
#!/usr/bin/with-contenv sh

exec /opt/kibana/bin/kibana
