# compose/django/s6/celerybeat/finish
#!/bin/sh

rm -f /tmp/celerybeat.pid

s6-svscanctl -t /etc/s6
