# compose/django/s6/celerybeat/run
#!/usr/bin/with-contenv sh
cd /app
exec celery -A velo.taskapp beat -l INFO --pidfile=/tmp/celerybeat.pid -S django  --uid=9005
