# compose/django/s6/celerybeat/run
#!/bin/sh
cd /app
exec celery -A velo.taskapp beat -l INFO --pidfile=
