# compose/django/s6/celerybeat/run
#!/bin/sh

exec cd /app && celery -A velo.taskapp beat -l INFO
